# deploy/aliyun/systemd/omnigraph-kol-classify.service
[Unit]
Description=OmniGraph daily KOL Layer-1 classify (5 topics)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/OmniGraph-Vault
EnvironmentFile=/root/.hermes/.env
ExecStart=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/batch_classify_kol.py --topic Agent --topic LLM --topic RAG --topic NLP --topic CV --min-depth 2 --days-back 1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
