# deploy/aliyun/systemd/omnigraph-vertex-proxy.service
[Unit]
Description=OmniGraph SSH SOCKS5 egress proxy: Aliyun -> GCP SG (Vertex/Google API, #75 B-mitigation)
After=network-online.target
Wants=network-online.target
# Disable start-limit so proxy always self-heals after transient SSH failures.
# Restart=always + RestartSec=30 provides natural backoff without a burst ceiling.
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
# -D opens SOCKS5 server on Aliyun loopback port 18080
# Python processes route Google/Vertex traffic via GCP Singapore VM (35.198.243.36)
# The embedding code uses httpx.AsyncClient(proxy="socks5h://127.0.0.1:18080") injected via
# HttpOptions(httpx_async_client=...) to bypass aiohttp's lack of native SOCKS5 support.
# Rollback: disable service + remove OMNIGRAPH_EMBED_PROXY from /root/.hermes/.env
# when IT confirms wg-gcp-sg ACK NetworkPolicy is fixed.
ExecStart=/usr/bin/ssh -N -T \
  -i /root/.ssh/id_ed25519 \
  -o ConnectTimeout=15 \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o BatchMode=yes \
  -o StrictHostKeyChecking=accept-new \
  -D 127.0.0.1:18080 \
  huhai_orion@35.198.243.36
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
