# deploy/aliyun/systemd/omnigraph-rss-fetch.timer
[Unit]
Description=OmniGraph daily RSS fetch timer
Requires=omnigraph-rss-fetch.service

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
