# deploy/aliyun/systemd/omnigraph-daily-digest.timer
[Unit]
Description=OmniGraph daily digest timer
Requires=omnigraph-daily-digest.service

[Timer]
OnCalendar=*-*-* 12:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
