# deploy/aliyun/systemd/omnigraph-kol-refresh.service
[Unit]
Description=OmniGraph KOL pre-scan cookie refresh (proactive — ssh-hermes CDP refresh before the daily scan)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
# Proactive refresh (KCA — "never block the scan"): ssh DIRECT to Hermes PC (NOT the retired Mac
# reverse tunnel) and run the refresh script SYNCHRONOUSLY (not detached — we want a fresh cookie
# BEFORE the daily scan fires). Bound the whole thing so a hung/unreachable Hermes never blocks:
#   - ConnectTimeout=20 caps the ssh dial
#   - the remote `timeout 240` caps the script (level-A is seconds; level-C QR-poll is ~5min but we
#     do NOT want pre-scan to block that long — if it needs a human QR scan, the script already
#     Telegram-notifies, and the scan will fail-then-OnFailure-retry as the belt-and-suspenders path)
#   - the || ... systemd-cat makes any failure visible in the journal, but DOES NOT fail the scan:
#     this unit is best-effort; the scan + its OnFailure alert remain the safety net.
# 2026-07-14: switched from `ssh hermes` (localhost:49221 → retired Mac tunnel, User hai, calling a
# non-existent run_refresh.sh — silently failing) to DIRECT Hermes PC (ohca.ddns.net:49221,
# User sztimhdd, /home/sztimhdd/OmniGraph-Vault). See docs/ALIYUN-WECHAT-REFRESH-OPS.md.
ExecStart=/bin/bash -c 'ssh -o BatchMode=yes -o ConnectTimeout=20 -o StrictHostKeyChecking=accept-new -p 49221 sztimhdd@ohca.ddns.net "cd ~/OmniGraph-Vault && timeout 240 python3 scripts/refresh_wechat_cookie.py >> ~/.hermes/kol-refresh.log 2>&1" || echo "kol-refresh: pre-scan ssh-hermes refresh failed (scan + OnFailure remain the safety net)" | systemd-cat -t kol-refresh'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
