# deploy/aliyun/systemd/omnigraph-daily-ingest.timer
[Unit]
Description=OmniGraph daily ingest 09:00 ADT (12:00 UTC) timer
Requires=omnigraph-daily-ingest.service

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
