# deploy/aliyun/systemd/omnigraph-rss-fetch.service
[Unit]
Description=OmniGraph daily RSS fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/OmniGraph-Vault
EnvironmentFile=/root/.hermes/.env
ExecStart=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/enrichment/rss_fetch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
