# deploy/aliyun/systemd/omnigraph-kol-classify.timer
[Unit]
Description=OmniGraph daily KOL Layer-1 classify (5 topics) timer
Requires=omnigraph-kol-classify.service

[Timer]
OnCalendar=*-*-* 11:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
