# deploy/aliyun/systemd/omnigraph-daily-ingest.service
[Unit]
Description=OmniGraph daily ingest 09:00 ADT (12:00 UTC)
After=network-online.target omnigraph-vertex-proxy.service
Wants=network-online.target omnigraph-vertex-proxy.service
Conflicts=omnigraph-afternoon-ingest.service omnigraph-evening-ingest.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/OmniGraph-Vault
EnvironmentFile=/root/.hermes/.env
TimeoutStartSec=300
MemoryHigh=3G
MemoryMax=4G
OOMScoreAdjust=500
Restart=on-failure
RestartSec=10min
ExecStartPre=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/scripts/cleanup_stuck_docs.py --all-failed
ExecStart=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/batch_ingest_from_spider.py --from-db --max-articles 10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
