# deploy/aliyun/systemd/omnigraph-kol-scan-batch@3.timer
[Unit]
Description=OmniGraph KOL scan batch 3 timer (11:00 UTC, staleness-partitioned)

[Timer]
OnCalendar=*-*-* 11:00:00 UTC
Persistent=true
Unit=omnigraph-kol-scan-batch@3.service

[Install]
WantedBy=timers.target
