# deploy/aliyun/systemd/omnigraph-kol-zombie-cleanup.timer
[Unit]
Description=OmniGraph kol zombie cleanup (pre-scan health check) timer
Requires=omnigraph-kol-zombie-cleanup.service

[Timer]
OnCalendar=*-*-* 10:55:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
