# deploy/aliyun/systemd/omnigraph-evening-ingest.timer
[Unit]
Description=OmniGraph evening ingest 21:00 ADT (00:00 UTC next day) timer

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
