# deploy/aliyun/systemd/omnigraph-evening-ingest.service
[Unit]
Description=OmniGraph evening ingest 21:00 ADT (00:00 UTC next day)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/OmniGraph-Vault
EnvironmentFile=/root/.hermes/.env
ExecStartPre=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/scripts/cleanup_stuck_docs.py --all-failed
ExecStart=/root/OmniGraph-Vault/venv-aim1/bin/python /root/OmniGraph-Vault/batch_ingest_from_spider.py --from-db --max-articles 5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
