# deploy/aliyun/systemd/omnigraph-daily-ingest.timer
# 2026-06-24 freq-up: consolidated sole ingest runner, every 2h (was 3 units 8h apart).
# afternoon + evening timers disabled (systemctl disable --now) — only this timer fires now.
[Unit]
Description=OmniGraph ingest every 2h (consolidated sole runner)

[Timer]
OnCalendar=*-*-* 00/2:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
