# deploy/aliyun/systemd/omnigraph-rss-rescrape.timer
[Unit]
Description=OmniGraph daily RSS body rescrape timer
Requires=omnigraph-rss-rescrape.service

[Timer]
OnCalendar=*-*-* 09:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
