# deploy/aliyun/systemd/omnigraph-kol-scan.timer
[Unit]
Description=OmniGraph daily KOL scan (WeChat MP API) timer
Requires=omnigraph-kol-scan.service

[Timer]
OnCalendar=*-*-* 11:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
