# deploy/aliyun/systemd/omnigraph-daily-ingest.service.d/override.conf
# Wave 3 T9 cutover 2026-06-02 — switch ingest to Qdrant vector storage (#25 close path)
# LightRAG 1.4.16 mandates QDRANT_URL when vector_storage="QdrantVectorDBStorage"
[Service]
Environment="OMNIGRAPH_VECTOR_STORAGE=qdrant"
Environment="QDRANT_URL=http://127.0.0.1:6333"

# 2026-06-10 hot-fix: bound asyncio post-completion hang at systemd layer (ISSUES #45)
# Root cause: lib/vision_tracking.drain_vision_tasks + LightRAG async pool may not close cleanly
# RuntimeMaxSec=3600 = 1h hard wall (2026-06-24 freq-up: was 3h; must stay < the 2h timer interval so a hung run is killed before next fire)
# TimeoutStopSec=300 = on Conflicts= SIGTERM, escalate to SIGKILL after 5min if Python doesn't exit
TimeoutStopSec=300
RuntimeMaxSec=3600
